# web_interface/adaptive-interpreter.service
[Unit]
Description=🧬💜 Adaptive Interpreter - Beautiful Variant Analysis Interface
After=network.target
Wants=network.target

[Service]
Type=simple
User=thereny
Group=thereny
WorkingDirectory=/home/Ace/AdaptiveInterpreter/web_interface
Environment=PYTHONPATH=/home/Ace/AdaptiveInterpreter
ExecStart=/usr/bin/python3 /home/Ace/AdaptiveInterpreter/web_interface/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/Ace/AdaptiveInterpreter

[Install]
WantedBy=multi-user.target
